# src/home/_blocklist/bin/blocklist.cron
#!/bin/sh
set -e

##
# environment
export HOME="/home/_blocklist"
export PATH="/usr/bin:/usr/local/bin:${PATH}"

##
# setup
. /home/_blocklist/.local/functions/date.sh
blocklist=$(which blocklist)
doas=$(which doas)
libexec="/home/_blocklist/.local/libexec"
path="/home/_blocklist/blocklists/$(today)"

##
# main
${blocklist} download
${blocklist} echo > ${path}
${doas} -u root ${libexec}/blocklist-copy
${doas} -u root ${libexec}/blocklist-reload
