# share/blocks+cron/crontab
0 0 * * * /home/_blocks/.local/bin/blocks+cron
